# quadlets/holland-backend.container
[Unit]
Description=Holland Backend (FastAPI)

[Container]
ContainerName=holland-backend
Image=ghcr.io/sebag90/holland-backend:latest
Volume=holland-uploads.volume:/data/uploads
# Mount your htpasswd file (bcrypt, generated with `htpasswd -B`):
Volume=/home/seba/containers_data/holland/htpasswd:/etc/holland/htpasswd:ro,Z

Environment=UPLOAD_DIR=/data/uploads
Environment=HTPASSWD_FILE=/etc/holland/htpasswd
EnvironmentFile=/home/seba/containers_data/holland/backend.env

Pod=holland.pod

[Service]
Restart=always

[Install]
WantedBy=default.target
